# webwatch.service
[Unit]
Description=Webwatch automatic web page tester
After=network.target

[Service]
ExecStart=/usr/local/bin/webwatch

