# webwatch.timer
[Unit]
Description=Webwatch automatic web page tester

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target

